# examples/service/feique.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Static preview of the launchd user agent generated by:
    feique service print --platform darwin

  DO NOT USE THIS FILE AS-IS. Replace <PLACEHOLDERS> with real paths,
  or (recommended) regenerate it on your own machine:

    feique service install --platform darwin \
        --config <USER_HOME>/.feique/config.toml

  See ./README.md for details.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>feique</string>
  <key>ProgramArguments</key>
  <array>
    <string><NODE_BIN></string>
    <string><FEIQUE_CLI></string>
    <string>serve</string>
    <string>--config</string>
    <string><CONFIG_PATH></string>
  </array>
  <key>WorkingDirectory</key>
  <string><WORKING_DIR></string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string><LOG_DIR>/feique.out.log</string>
  <key>StandardErrorPath</key>
  <string><LOG_DIR>/feique.err.log</string>
  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
